# echo@.service
[Unit]
Description=Defining routine for echo service

[Service]
ExecStart=-/bin/cat -n
StandardInput=socket
